# config/etc/supervisor.d/httpd.conf
#Donar-li nom al servei
[program:httpd]
#Executa el servei httpd en segon pla
command=/usr/sbin/httpd -DFOREGROUND
