# config/etc/supervisor.d/minidlnad.conf
# Donar-li nom al servei
[program:minidlnad]
# Executar el servei Minidlna en segon pla
command=/sbin/minidlnad -d
